# config/supervisord.conf
[supervisord]
nodaemon=true
[program:php-fpm]
autorestart = true
command = /usr/local/php/sbin/php-fpm
